# cluster/gitops/gitops.timer
[Unit]
Description=Rodar o sincronizador GitOps a cada 1 minuto

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
