# cluster/autoscaler/autoscaler.timer
[Unit]
Description=Rodar autoscaler a cada minuto

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
